# config/messenger-worker.ini
[program:messenger-consume]
command=php /home/krish.kalaria@simform.dom/Desktop/LMS-2/start/bin/console messenger:consume async_priority_high async --time-limit=3600
user=krish
numprocs=2
startsecs=0
autostart=true
autorestart=true
startretries=10
process_name=%(program_name)s_%(process_num)02d